# deploy/bitget-scan.timer
[Unit]
Description=Bitget paper scan cadence (5 min ; large marge sous les limites API Bitget)

[Timer]
# 5 min : ~30-100x sous l'API (20/s par endpoint, 100/s agrege) meme sur l'univers
# dynamique (15 symboles). Notifications Telegram deportees -> bitget-notify.timer.
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
