# deploy/bitget-logrotate.timer
[Unit]
Description=Bitget rotation des journaux (quotidien)

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
