# deploy/bitget-security-audit.service
[Unit]
Description=Bitget audit securite quotidien (VERDICT SAFE attendu)

[Service]
Type=oneshot
WorkingDirectory=/root/bitget_termux_repo
ExecStart=/usr/bin/python3 security_agent.py
User=root
Nice=15
