# deploy/bitget-watchdog.timer
[Unit]
Description=Bitget watchdog cadence (5 min ; alerte immediate sur DOWN/STALE)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=20s
Persistent=true

[Install]
WantedBy=timers.target
